# files/mailscanner-rootdir/etc/cron.d/sync
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=""
48 * * * * root sync; echo 3 > /proc/sys/vm/drop_caches
38 23 * * * root fstrim /
